# .config/autostart/noding.desktop
[Desktop Entry]
Type=Application
Name=Stop the Ding!
Exec=xset -b b off
